# com.lifeos.tldr_ingest.example.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.lifeos.tldr_ingest</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/path/to/LifeOS/.venv/bin/python</string>
        <string>/path/to/LifeOS/scripts/auto_tldr.py</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/path/to/LifeOS</string>
    
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/path/to/LifeOS/logs/tldr_launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/path/to/LifeOS/logs/tldr_launchd_err.log</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/path/to/LifeOS/.venv/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>PYTHONPATH</key>
        <string>/path/to/LifeOS</string>
    </dict>
</dict>
</plist>
